# tools/http-port-manager/launchd/com.user.http-port-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.http-port-manager</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__SERVER_PY__</string>
    <string>--config</string>
    <string>__CONFIG_JSON__</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__APP_DIR__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/http-port-manager.out.log</string>

  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/http-port-manager.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
